# etc/supervisord.d/owncloud.ini
[program:owncloud]
command = consul-template -consul localhost:8500 -wait 5s -template "/etc/consul-templates/owncloud/owncloud.conf.ctmpl:/etc/nginx/conf.d/owncloud.conf:nginx -s reload"
redirect_stderr=true
stdout_logfile = /var/log/supervisor/%(program_name)s.log
