# etc/supervisord.d/owncloud.ini
[program:owncloud]
command=/opt/qnib/owncloud/bin/start.sh
redirect_stderr=true
stdout_logfile = /var/log/supervisor/%(program_name)s.log
