# etc/supervisord.d/php-fpm.ini
[program:php-fpm]
command=/usr/bin/php-fpm --nodaemonize
redirect_stderr=true
stdout_logfile=/var/log/supervisor/%(program_name)s.log
